# debian/grase-conf-openvpn.cron.daily
#!/bin/bash

# Check for signed certificate from server and install
/usr/share/grase/scripts/retrieve_key_certificate

