# debian/grase-conf-openvpn.cron.hourly
#!/bin/bash

if [ -x "/usr/share/grase/scripts/get_signed_cert" ]; then
	/usr/share/grase/scripts/get_signed_cert
fi
