# debian/grase-conf-openvpn.cron.hourly
#!/bin/bash

if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d openvpn status grase 2>&1 > /dev/null
	ret=$?
else
	/etc/init.d/openvpn status grase > /dev/null
	ret=$?
fi

if [ $ret -eq 1 ]; then

    if which invoke-rc.d >/dev/null 2>&1; then
	    invoke-rc.d openvpn restart grase
    else
	    /etc/init.d/openvpn restart grase
    fi
    
fi    
